# iptables.upstart
description "Load /etc/sysconfig/iptables pre-defined rules"

start on pld.network-started
stop on pld.network-stopped

env IPTABLES_CONFIG=/etc/sysconfig/iptables

pre-start script
	# fail early if config missing
	test -f $IPTABLES_CONFIG || /sbin/initctl stop
end script

script
	exec /etc/rc.d/init.d/iptables load
end script

pre-stop script
	exec /etc/rc.d/init.d/iptables clear
end script
